# user/base/inittab
::sysinit:/bin/splash
::respawn:/bin/desktop
::respawn:/etc/init.d/network-service
::respawn:-/bin/sh
