# user/base/inittab
::respawn:/bin/console-session
::respawn:/etc/init.d/network-service
::respawn:-/bin/sh
