# user/base/inittab
::once:/etc/init.d/graphical-session /bin/compositor
::once:/etc/init.d/graphical-session /bin/lite-ui --desktop
::respawn:/etc/init.d/network-service
::respawn:-/bin/sh
