# user/inittab
::respawn:/etc/init.d/network-service
::askfirst:-/bin/sh
