# user/inittab
::respawn:/etc/init.d/network-service
::respawn:/bin/liteos-terminal
::respawn:-/bin/sh
