# user/inittab
::once:/bin/display-session-guard
::respawn:/bin/liteos-terminal
::respawn:/etc/init.d/network-service
::respawn:-/bin/sh
